# 1046/openvswitch-filelist
/.
/usr
/usr/bin
/usr/bin/ovs-appctl
/usr/bin/ovs-docker
/usr/bin/ovs-ofctl
/usr/bin/ovs-parse-backtrace
/usr/bin/ovs-pki
/usr/bin/ovsdb-client
/usr/sbin
/usr/sbin/ovs-bugtool
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/ovs-appctl-bashcomp.bash
/usr/share/doc
/usr/share/doc/openvswitch-common
/usr/share/doc/openvswitch-common/changelog.Debian.gz
/usr/share/doc/openvswitch-common/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ovsdb-client.1.gz
/usr/share/man/man1/ovsdb-tool.1.gz
/usr/share/man/man8
/usr/share/man/man8/ovs-appctl.8.gz
/usr/share/man/man8/ovs-bugtool.8.gz
/usr/share/man/man8/ovs-ofctl.8.gz
/usr/share/man/man8/ovs-parse-backtrace.8.gz
/usr/share/man/man8/ovs-pki.8.gz
/usr/share/openvswitch
/usr/share/openvswitch/bugtool-plugins
/usr/share/openvswitch/bugtool-plugins/kernel-info
/usr/share/openvswitch/bugtool-plugins/kernel-info/openvswitch.xml
/usr/share/openvswitch/bugtool-plugins/network-status
/usr/share/openvswitch/bugtool-plugins/network-status/openvswitch.xml
/usr/share/openvswitch/bugtool-plugins/network-status/ovn.xml
/usr/share/openvswitch/bugtool-plugins/system-configuration
/usr/share/openvswitch/bugtool-plugins/system-configuration/openvswitch.xml
/usr/share/openvswitch/bugtool-plugins/system-configuration.xml
/usr/share/openvswitch/bugtool-plugins/system-logs
/usr/share/openvswitch/bugtool-plugins/system-logs/openvswitch.xml
/usr/share/openvswitch/scripts
/usr/share/openvswitch/scripts/ovs-bugtool-daemons-ver
/usr/share/openvswitch/scripts/ovs-bugtool-fdb-show
/usr/share/openvswitch/scripts/ovs-bugtool-get-dpdk-nic-numa
/usr/share/openvswitch/scripts/ovs-bugtool-ovs-appctl-dpif
/usr/share/openvswitch/scripts/ovs-bugtool-ovs-bridge-datapath-type
/usr/share/openvswitch/scripts/ovs-bugtool-ovs-ofctl-loop-over-bridges
/usr/share/openvswitch/scripts/ovs-bugtool-ovs-vswitchd-threads-affinity
/usr/share/openvswitch/scripts/ovs-bugtool-qos-configs
/usr/share/openvswitch/scripts/ovs-bugtool-tc-class-show
/usr/share/openvswitch/scripts/ovs-lib
/var
/var/log
/var/log/openvswitch
/.
/etc
/etc/default
/etc/default/openvswitch-switch
/etc/init.d
/etc/init.d/openvswitch-switch
/etc/logrotate.d
/etc/logrotate.d/openvswitch-switch
/etc/network
/etc/network/if-post-down.d
/etc/network/if-pre-up.d
/etc/openvswitch
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/openvswitch-switch.service
/lib/systemd/system/ovs-vswitchd.service
/lib/systemd/system/ovsdb-server.service
/usr
/usr/bin
/usr/bin/ovs-dpctl
/usr/bin/ovs-dpctl-top
/usr/bin/ovs-pcap
/usr/bin/ovs-tcpdump
/usr/bin/ovs-tcpundump
/usr/bin/ovs-vlan-test
/usr/bin/ovs-vsctl
/usr/bin/ovsdb-tool
/usr/lib
/usr/lib/openvswitch-switch
/usr/lib/openvswitch-switch/ovs-vswitchd
/usr/sbin
/usr/sbin/ovsdb-server
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/ovs-vsctl-bashcomp.bash
/usr/share/doc
/usr/share/doc/openvswitch-switch
/usr/share/doc/openvswitch-switch/README.Debian.gz
/usr/share/doc/openvswitch-switch/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ovs-pcap.1.gz
/usr/share/man/man1/ovs-tcpundump.1.gz
/usr/share/man/man1/ovsdb-server.1.gz
/usr/share/man/man5
/usr/share/man/man5/ovs-vswitchd.conf.db.5.gz
/usr/share/man/man8
/usr/share/man/man8/ovs-ctl.8.gz
/usr/share/man/man8/ovs-dpctl-top.8.gz
/usr/share/man/man8/ovs-dpctl.8.gz
/usr/share/man/man8/ovs-tcpdump.8.gz
/usr/share/man/man8/ovs-vsctl.8.gz
/usr/share/man/man8/ovs-vswitchd.8.gz
/usr/share/openvswitch
/usr/share/openvswitch/scripts
/usr/share/openvswitch/scripts/ifupdown.sh
/usr/share/openvswitch/scripts/ovs-check-dead-ifs
/usr/share/openvswitch/scripts/ovs-ctl
/usr/share/openvswitch/scripts/ovs-save
/usr/share/openvswitch/switch
/usr/share/openvswitch/switch/default.template
/usr/share/openvswitch/vswitch.ovsschema
/etc/network/if-post-down.d/openvswitch
/etc/network/if-pre-up.d/openvswitch
/usr/share/doc/openvswitch-switch/changelog.Debian.gz
/usr/local/bin/ovs-l3ping
/usr/local/bin/ovs-test
/usr/local/bin/ovs-testcontroller
/usr/local/sbin/ovs-vlan-bug-workaround
/usr/local/sbin/ovs-vswitchd
